# src/udev/10-kissmodem.rules
SUBSYSTEM=="tty", SUBSYSTEMS=="usb", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="7523", SYMLINK+="kiss", RUN+="/bin/systemctl restart kiss-modem"
